# key_tempo_bot_service.service
[Unit]
Description=key_tempo_bot

[Service]
ExecStart=/usr/bin/python3.9 ~/key_tempo_bot/finder.py

